# etc/profile.d/wayland-portals.sh
export GTK_USE_PORTAL=1
export MOZ_ENABLE_WAYLAND=1
# XDG_CURRENT_DESKTOP non va impostato qui:
# viene esportato da niri-session via `systemctl --user import-environment`
