# etc/systemd/system/helium-setup.service
[Unit]
Description=Install Helium Flatpak from pre-downloaded bundle
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Check if the flatpak is already installed; if not, install the bundle
ExecStartPre=/usr/bin/sh -c 'if ! flatpak info net.imput.helium >/dev/null 2>&1; then exit 0; else exit 1; fi'
ExecStart=/usr/bin/flatpak install --system -y /usr/share/helium/helium.flatpak
# Crea il link simbolico per ogni utente esistente per permettere a Helium (sandbox) di trovarsi
ExecStartPost=/usr/bin/bash -c 'for user_home in /var/home/*; do \
    if [ -d "$user_home" ]; then \
        mkdir -p "$user_home/.local/share/applications"; \
        ln -sf /var/lib/flatpak/exports/share/applications/net.imput.helium.desktop "$user_home/.local/share/applications/net.imput.helium.desktop"; \
        chown -R $(basename "$user_home"):$(basename "$user_home") "$user_home/.local/share/applications"; \
    fi; \
done'
# Optional: clean up the bundle after installation to save space in the persistent /var if needed
# ExecStartPost=/usr/bin/rm -f /usr/share/helium/helium.flatpak
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
