# etc/systemd/system/helium-setup.service
[Unit]
Description=Install Helium Flatpak from pre-downloaded bundle
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Check if the flatpak is already installed; install if missing
ExecStartPre=/usr/bin/sh -c 'if ! flatpak info net.imput.helium >/dev/null 2>&1; then exit 0; else exit 0; fi'
ExecStart=/usr/bin/flatpak install --system -y /usr/share/helium/helium.flatpak
# Imposta Helium come opzione disponibile (ma non forzata come predefinita)
ExecStartPost=/usr/bin/bash -c 'for user_home in /var/home/*; do \
    if [ -d "$user_home" ]; then \
        MIME_FILE="$user_home/.config/mimeapps.list"; \
        mkdir -p "$(dirname "$MIME_FILE")"; \
        if [ ! -f "$MIME_FILE" ]; then echo "[Default Applications]" > "$MIME_FILE"; fi; \
        if ! grep -q "\[Added Associations\]" "$MIME_FILE"; then echo "[Added Associations]" >> "$MIME_FILE"; fi; \
        if ! grep -q "net.imput.helium.desktop" "$MIME_FILE"; then \
            sed -i "/\[Added Associations\]/a text/html=net.imput.helium.desktop;x-scheme-handler/http=net.imput.helium.desktop;x-scheme-handler/https=net.imput.helium.desktop;" "$MIME_FILE"; \
        fi; \
        ln -sf /var/lib/flatpak/exports/share/applications/net.imput.helium.desktop "$user_home/.local/share/applications/helium.desktop"; \
        chown $(basename "$user_home"):$(basename "$user_home") "$user_home/.local/share/applications/helium.desktop"; \
    fi; \
done'
# Optional: clean up the bundle after installation to save space in the persistent /var if needed
# ExecStartPost=/usr/bin/rm -f /usr/share/helium/helium.flatpak
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
